# qr3_64.conf
[program:qr3_64]
command=python3 /home/curator3/qr3_64/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_64/test.err.log
stdout_logfile=/home/curator3/qr3_64/test.err.log
